# smf/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='application-fs-caching-server'>
	<service name='application/fs-caching-server' type='service' version='0'>
		<create_default_instance enabled='true'/>
		<dependency name='dep0' grouping='require_all' restart_on='error' type='service'>
			<service_fmri value='svc:/milestone/multi-user:default'/>
		</dependency>
		<method_context working_directory='/tmp'>
			<method_credential user='nobody' group='other'/>
			<method_environment>
				<envvar name='FS_CACHE_DIR' value='/var/tmp/cache'/> <!-- the cache dir -->
				<envvar name='FS_CACHE_URL' value='http://pkgsrc.joyent.com'/> <!-- the URL you want to proxy -->
				<envvar name='FS_CACHE_HOST' value='0.0.0.0'/>
				<envvar name='FS_CACHE_PORT' value='8080'/>
				<envvar name='FS_CACHE_REGEX' value='\.(gif|png|jpg|tar|tgz|tar\.gz)$'/> <!-- leave blank for default -->
				<envvar name='FS_CACHE_DEBUG' value=''/> <!-- set to anything to enable -->
			</method_environment>
		</method_context>
		<exec_method name='start' type='method' exec='fs-caching-server &amp;' timeout_seconds='10'/>
		<exec_method name='stop' type='method' exec=':kill' timeout_seconds='30'/>
		<template>
			<common_name>
				<loctext xml:lang='C'>FS Caching Server</loctext>
			</common_name>
		</template>
	</service>
</service_bundle>
